# com.example.dailyreport.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.example.dailyreport</string>
    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <!-- UPDATE: path to where you cloned agent-timecard -->
      <string>/Users/YOURUSERNAME/code/agent-timecard/daily_report.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
      <key>Hour</key>
      <integer>8</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <!-- UPDATE: must include uv, claude, and node in PATH -->
      <string>/usr/local/bin:/usr/bin:/bin</string>
      <key>HOME</key>
      <!-- UPDATE: your home directory -->
      <string>/Users/YOURUSERNAME</string>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/dailyreport.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/dailyreport.err</string>
  </dict>
</plist>
